# systemd/defendos-healthcheck.timer
[Unit]
Description=Run DefendOS scheduled healthcheck daily at 09:00

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true
Unit=defendos-healthcheck.service

[Install]
WantedBy=timers.target
